# config/crontab/crontab
* * * * * cd /var/www/html/Pure_BE && php artisan schedule:run >> /dev/null 2>&1